# deploy/coastal-dashboard.service
[Unit]
Description=Coastal Emergency Agent (dashboard)
After=network-online.target coastal-agent.service
Wants=network-online.target

[Service]
Type=simple
User=coastal
Group=coastal
WorkingDirectory=/home/coastal/Coastal_Emergency_Agent_EGU26
EnvironmentFile=/opt/coastal-agent.env
ExecStart=/home/coastal/Coastal_Emergency_Agent_EGU26/.venv/bin/uvicorn dashboard.api:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
